# etc/systemd/system/audio-speakers.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/audio-speakers.sh

[Install]
WantedBy=default.target

